# esxi-on-kvm/config/60-vnet.rules
KERNEL=="ovs*", RUN+="/sbin/ip link set %k txqueuelen 5000"
KERNEL=="esxi*", RUN+="/sbin/ip link set %k txqueuelen 5000"
KERNEL=="eth0", RUN+="/sbin/ip link set %k txqueuelen 5000" 
KERNEL=="lo", RUN+="/sbin/ip link set %k txqueuelen 5000"